# deploy/cleanup/prestashop-log-cleanup.cron.daily
#!/usr/bin/env bash
# Wrapper fuer /etc/cron.daily: datierte PrestaShop-Logs bereinigen.
set -euo pipefail

/usr/local/sbin/prestashop-dated-logs-cleanup /var/www/wow9.internaut.ch/html/var/logs 30
