# deploy/cleanup/prestashop-log-cleanup.cron.daily
#!/usr/bin/env bash
# Wrapper fuer /etc/cron.daily: datierte PrestaShop-Logs bereinigen.
set -euo pipefail

/usr/local/sbin/prestashop-dated-logs-cleanup /path/to/prestashop/var/logs 30
